# rootfs/etc/s6-overlay/s6-rc.d/loki/run
#!/command/with-contenv bash
exec s6-setuidgid unhealthy /usr/local/bin/loki \
  -config.file=/etc/loki/loki.yml \
  -config.expand-env=true
