# rootfs/etc/s6-overlay/s6-rc.d/init-prometheus-web-config/up
/command/with-contenv /etc/cont-init.d/03-prometheus-web-config
