# rootfs/etc/s6-overlay/s6-rc.d/bgw-syslog-forwarder/run
#!/command/with-contenv bash
exec /usr/local/bin/bgw-syslog-forwarder
