# rootfs/etc/s6-overlay/s6-rc.d/init-textfile-dir/up
/command/with-contenv /etc/cont-init.d/02-textfile-dir
