# rootfs/etc/s6-overlay/s6-rc.d/unifi-protect-exporter/run
#!/command/with-contenv bash
# unifi-protect-exporter — poll the UniFi Protect Integration API for UP-Sense
# environmental sensors and expose them on 127.0.0.1:9688. Configure the NVR via
# UNIFI_PROTECT_HOST + UNIFI_PROTECT_API_KEY. Read-only; idles with
# unifi_sensor_scrape_success=0 when unset or unreachable.
exec /usr/local/bin/unifi-protect-exporter
