# rootfs/etc/s6-overlay/s6-rc.d/prometheus/run
#!/command/with-contenv bash
exec s6-setuidgid unhealthy /usr/local/bin/prometheus \
  --config.file=/etc/prometheus/prometheus.yml \
  --storage.tsdb.path=/var/lib/prometheus \
  --storage.tsdb.retention.time=${PROMETHEUS_RETENTION:-30d} \
  --web.listen-address=127.0.0.1:9090 \
  --web.enable-lifecycle
